# templates/init.conf.erb
# generated by puppet, changes will be overwritten

start on <%= @start_on %>
stop on <%= @stop_on %>

<% if @expect -%>
expect <%= @expect %>
<% end -%>
<% chpst_args = nil -%>
<% if @chdir -%>
chdir <%= @chdir %>
<% end -%>
<% if @chroot -%>
chroot <%= @chroot %>
<% end -%>
<% if @user and scope.lookupvar("upstart::params::setuid_supported") %>
setuid <%= @user %>
<% elsif @user %>
<% chpst_args += " -u #{@user}" -%>
<% end -%>
<% if @group and scope.lookupvar("upstart::params::setgid_supported") %>
setgid  <%= @group %>
<% elsif @group -%>
<% chpst_args += " -g #{@group}" -%>
<% end -%>
<% if @nice -%>
nice <%= @nice %>
<% end -%>

<% @env.each do |k,v| -%>
env <%= "#{k}=\"#{v.gsub('\"', '\\\"')}\"" %>
<% end -%>

<% @limits.each do |k,v| -%>
limit <%= k %> <%= v %> <%= v %>
<% end -%>

<% if @exec -%>
exec <%= chpst_args ? "chpst.py #{chpst_args} -- #{@exec}" : @exec %>
<% elsif @script %>
script
	<%= @script %>
end script
<% end %>
