# templates/init.conf.erb
# generated by puppet, changes will be overwritten

start on <%= $start_on %>
stop on <%= $stop_on %>

<% chpst_args = "" -%>
<% if @chdir -%>
chdir <%= @chdir %>
<% end -%>
<% if @user and scope.lookupvar("upstart::params::setuid_supported") %>
setuid <%= @user %>
<% elsif @user %>
<% chpst_args += "-u #{@user}" -%>
<% end -%>
<% if @group and scope.lookupvar("upstart::params::setgid_supported") %>
<% chpst_args += "-g #{@group}" -%>
<% end -%>

<% @env.each do |k,v| -%>
<%= "#{k}=#{v}" %>
<% end -%>

<% if @exec -%>
exec <%= chpst_args ? "chpst.py #{chpst_args} -- #{@exec}" : @exec %>
<% elsif @script %>
script
	<%= @script %>
end script
<% end %>